# files/default/heroprotocol.cron
#!/bin/sh

cd /opt/heroprotocol
git pull origin master