# files/default/heroprotocol.cron
#!/bin/sh

cd /opt/heroprotocol
git pull
chmod 755 /opt/heroprotocol/heroprotocol.py