# files/default/laravel.cron
* * * * * www-data php /var/www/hotsapi/artisan schedule:run >> /dev/null 2>&1